# apps/copy-engine/deploy/supervisord.conf
; AlgoSphere copy-engine — Supervisor config (single-VPS alternative).
;
; Use this when running all workers under one supervisor on a worker VPS
; (instead of one container per worker via docker-compose). Redis is assumed
; to be reachable via REDIS_URL (run it separately or on the same host).
; Logs go to stdout/stderr so a process manager / journald / Loki agent can
; collect them. autorestart + startretries give crash recovery.
;
;   supervisord -c deploy/supervisord.conf
;   supervisorctl -c deploy/supervisord.conf status

[supervisord]
nodaemon=true
logfile=/var/log/copy-engine/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
loglevel=info

[unix_http_server]
file=/tmp/copy-engine-supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/copy-engine-supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:orchestrator]
command=python orchestrator.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:executor]
command=python executor.py
directory=/app
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:reconciler]
command=python reconciler.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:coach]
command=python coach.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[group:copy-engine]
programs=orchestrator,executor,reconciler,coach
